# actividad4/testserv.service
[Unit]
Description=My Script

[Service]
Type=simple
ExecStart=/media/script/test.sh
Restart=always

[Install]
WantedBy=multi-user.target
